# services/com.autoforge.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  AutoForge launchd service (macOS)

  Install:
    cp services/com.autoforge.daemon.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.autoforge.daemon.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.autoforge.daemon.plist
    rm ~/Library/LaunchAgents/com.autoforge.daemon.plist

  Logs:
    tail -f /tmp/autoforge-daemon.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.autoforge.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "$HOME/autoforge" &amp;&amp; .venv/bin/python forge.py daemon start</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/autoforge-daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/autoforge-daemon.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
